# apps/cpoe_cli/packaging/linux/systemd/cpoe-ibus.service
[Unit]
Description=CPOE IBus Input Method Engine
Documentation=man:writerslogic(1)
Documentation=https://github.com/writerslogic/writerslogic
After=ibus.service
BindsTo=graphical-session.target
PartOf=graphical-session.target

# Require IBus to be running
Requires=ibus.socket
After=ibus.socket

[Service]
Type=dbus
BusName=org.freedesktop.IBus.CPOE
ExecStart=/usr/bin/writerslogic-ibus
ExecReload=/bin/kill -HUP $MAINPID
Restart=on-failure
RestartSec=3
TimeoutStartSec=10
TimeoutStopSec=10

# Environment
Environment=CPOE_IBUS_ENGINE=1

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=writerslogic-ibus

# Security
NoNewPrivileges=yes
PrivateTmp=yes

# The IBus engine needs access to X11/Wayland
# so we cannot fully sandbox it

[Install]
WantedBy=graphical-session.target
