# apps/cpoe_cli/packaging/linux/systemd/cpoe.socket
[Unit]
Description=CPOE Socket Activation
Documentation=man:cpoe(1)
PartOf=cpoe.service

[Socket]
# Unix socket for local communication
ListenStream=/run/cpoe/cpoe.sock
SocketMode=0660
SocketUser=cpoe
SocketGroup=cpoe

# Accept connections
Accept=no

# Buffer sizes
ReceiveBuffer=1M
SendBuffer=1M

# Backlog
Backlog=128

# Permissions
DirectoryMode=0750

[Install]
WantedBy=sockets.target
