# apps/cpoe_cli/packaging/linux/systemd/cpoe-user.service
[Unit]
Description=WritersProof authorship witnessing daemon (user session)
Documentation=man:writersproof-cli(1)
Documentation=https://github.com/writerslogic/writersproof-cli
After=graphical-session.target

[Service]
# Runs as the logged-in user inside the graphical session. This is the only
# model that works for focus tracking: it has DISPLAY/WAYLAND_DISPLAY, and it
# can read /dev/input when the user is in the `input` group. A *system* service
# cannot do either (PrivateDevices blocks /dev/input; a system user has no
# display), which is why no system unit is shipped.
Type=simple
ExecStart=/usr/bin/writersproof-cli start --foreground
ExecStop=/usr/bin/writersproof-cli stop
Restart=on-failure
RestartSec=5

StandardOutput=journal
StandardError=journal
SyslogIdentifier=writersproof-cli

[Install]
WantedBy=default.target
