# apps/cpoe_cli/packaging/linux/systemd/cpoe-user.service
[Unit]
Description=CPOE User Authorship Witnessing Service
Documentation=man:writerslogic(1)
Documentation=https://github.com/writerslogic/writerslogic
After=graphical-session.target

[Service]
Type=notify
ExecStart=/usr/bin/writerslogic --user
ExecReload=/bin/kill -HUP $MAINPID
Restart=on-failure
RestartSec=5
TimeoutStartSec=30
TimeoutStopSec=30

# Environment
Environment=CPOE_USER_MODE=1

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=writerslogic-user

# Security (user service, less restrictive)
NoNewPrivileges=yes
PrivateTmp=yes

[Install]
WantedBy=default.target
