# apps/cpoe_cli/packaging/linux/systemd/cpoe.service
[Unit]
Description=CPOE Cryptographic Authorship Witnessing Daemon
Documentation=man:writersproof-cli(1)
Documentation=https://github.com/writerslogic/writerslogic
After=network.target local-fs.target
Wants=writerslogic.socket

[Service]
Type=notify
User=writerslogic
Group=writerslogic

# Configuration
EnvironmentFile=-/etc/writerslogic/environment
ExecStart=/usr/bin/writersproof-cli --config /etc/writerslogic/config.toml
ExecReload=/bin/kill -HUP $MAINPID
Restart=on-failure
RestartSec=5
TimeoutStartSec=30
TimeoutStopSec=30
WatchdogSec=60

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=writerslogic

# Paths
RuntimeDirectory=writerslogic
RuntimeDirectoryMode=0750
StateDirectory=writerslogic
StateDirectoryMode=0750
LogsDirectory=writerslogic
LogsDirectoryMode=0750
ConfigurationDirectory=writerslogic
ConfigurationDirectoryMode=0750

# Security hardening
NoNewPrivileges=yes
PrivateTmp=yes
PrivateDevices=yes
ProtectSystem=strict
ProtectHome=read-only
ProtectKernelTunables=yes
ProtectKernelModules=yes
ProtectKernelLogs=yes
ProtectControlGroups=yes
ProtectHostname=yes
ProtectClock=yes
ProtectProc=invisible
ProcSubset=pid
RestrictAddressFamilies=AF_UNIX AF_INET AF_INET6
RestrictNamespaces=yes
RestrictRealtime=yes
RestrictSUIDSGID=yes
LockPersonality=yes
MemoryDenyWriteExecute=yes
RemoveIPC=yes
UMask=0077

# Capabilities
CapabilityBoundingSet=
AmbientCapabilities=

# Resource limits
LimitNOFILE=65536
LimitNPROC=64
LimitCORE=0

# cgroups resource control
CPUQuota=50%
MemoryMax=512M
TasksMax=64

# Sandboxing
SystemCallArchitectures=native
SystemCallFilter=@system-service
SystemCallFilter=~@privileged @resources @mount @swap @reboot @obsolete

[Install]
WantedBy=multi-user.target
